# packaging/rpm/forwarder.postinst
if [ $1 -eq 1 ] && [ -x /usr/bin/systemctl ]; then
        # Initial installation
        /usr/bin/systemctl --no-reload preset forwarder.service || :
fi
