# packaging/rpm/forwarder.prerm
if [ $1 -eq 0 ] && [ -x /usr/bin/systemctl ]; then
        # Package removal, not upgrade
        /usr/bin/systemctl --no-reload disable --now forwarder.service || :
fi
